# configs/etc/systemd/system/bootc-fetch-apply-updates.timer.d/override.conf
[Timer]
RandomizedDelaySec=10m
OnCalendar=daily
Persistent=true